# dummy.service
[Unit]
Description=Dummy Background Service
After=network.target

[Service]
ExecStart=/usr/local/bin/dummy.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
